# var/docker/supervisord/backend.ini
[program:backend]
directory=/app
command=npm run start:prod
autostart=true
autorestart=false
