# var/docker/supervisord/cron.ini
[program:cron]
directory=/app
command=npm run start:prod:cron
autostart=true
autorestart=false
