# var/docker/supervisord/frontend.ini
[program:frontend]
directory=/app
command=npm run start:prod:frontend
autostart=true
autorestart=false
