# var/docker/supervisord/workers.conf
[program:workers]
directory=/app/apps/workers
command=pnpm start
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
startsecs=10
