# var/docker/supervisord/workers.ini
[program:workers]
directory=/app
command=npm run start:prod:workers
autostart=true
autorestart=false
